# systemd/dns-syncer-sync.service
[Unit]
Description=DNS Syncer DNS update check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=dns-syncer
Group=dns-syncer
WorkingDirectory=/opt/dns-syncer
ExecStart=/opt/dns-syncer/.venv/bin/python -m app.cli sync-once
